# btbutton-mqtt@.service
[Unit]
Description=btbutton remote on %I

[Service]
ExecStart=/usr/local/bin/btbutton-mqtt.py /%I
Type=simple
